# cloudflare-tunnel.service
[Unit]
Description=Cloudflare Quick Tunnel for WiFi Controller Dashboard
After=network-online.target
Wants=network-online.target
Before=wifi-controller.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Desktop/WiFi-Switch
ExecStart=/home/pi/Desktop/WiFi-Switch/venv/bin/python3 /home/pi/Desktop/WiFi-Switch/cloudflare_tunnel.py

Restart=always
RestartSec=10
StartLimitIntervalSec=0

StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
