# wifi-controller.service
[Unit]
Description=WiFi Controller with Dashboard
After=network.target network-online.target cloudflare-tunnel.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Desktop/WiFi-Switch
ExecStart=/home/pi/Desktop/WiFi-Switch/venv/bin/python3 -u /home/pi/Desktop/WiFi-Switch/wifi_controller.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
